# interface/openwrt/files/etc/init.d/mahjong-mitm
#!/bin/sh /etc/rc.common
# mahjong-mitm — setup-period 热更 MITM（procd 服务）
#
# 启动 python -m mahjong_mitm，参数从 uci config (/etc/config/mahjong-mitm) 读取。
# 启用时加载 nftables DNS 劫持规则；停用时清除。

START=95
STOP=10
USE_PROCD=1

PROG_DIR="/usr/lib/mahjong-mitm"
PYTHON="/usr/bin/python3"

# 读取 uci，给不存在的字段兜底默认值
load_config() {
	config_load mahjong-mitm
	config_get enabled         main enabled         1
	config_get host_ip         main host_ip         192.168.8.1
	config_get ecs_ip          main ecs_ip          8.136.32.137
	config_get tls_port        main tls_port        443
	config_get dns_port        main dns_port        5353
	config_get dns_listen_host main dns_listen_host 0.0.0.0
}

apply_nft() {
	# fw4 会 include /etc/nftables.d/*.nft；reload 即生效
	[ -x /sbin/fw4 ] && fw4 reload >/dev/null 2>&1
}

apply_gxb_dnat() {
	# 解析 gxb-* 真实 CDN IP 并把它们的 :443 流量 DNAT 到本机。
	# 必须在 fw4 reload 之后执行（fw4 reload 会清空动态规则）。
	# 用途：手机 DNS 缓存了真实 IP 时，仍能被透明劫持到 MITM。
	local resolver="$PROG_DIR/resolve-gxb.sh"
	[ -x "$resolver" ] || return 0
	sh "$resolver" "$host_ip" 2>/tmp/mahjong-gxb-resolve.log | while read -r line; do
		[ -z "$line" ] && continue
		# 跳过 '#' 注释
		case "$line" in '#'*) continue ;; esac
		nft "$line" 2>>/tmp/mahjong-gxb-resolve.log || true
	done
}

clear_nft() {
	# 删掉规则文件再 reload，让 mahjong_mitm_dns 链消失（恢复正常 DNS）
	[ -x /sbin/fw4 ] && fw4 reload >/dev/null 2>&1
}

start_service() {
	load_config

	if [ "$enabled" != "1" ]; then
		echo "mahjong-mitm: disabled in uci config, not starting" >&2
		return 0
	fi

	# 装载 nftables DNS 劫持（含硬编码 DNS 重定向到本机 dns_port）
	apply_nft
	# 追加 gxb-* 真实 CDN IP 的 DNAT 规则（绕过手机 DNS 缓存）
	apply_gxb_dnat

	local args="--host-ip $host_ip --ecs-ip $ecs_ip --tls-port $tls_port --dns-port $dns_port --assets-dir $PROG_DIR/assets"
	[ -n "$dns_listen_host" ] && args="$args --dns-listen-host $dns_listen_host"

	# 清理旧实例残留的 443 占用，避免新实例 bind 失败。
	# OpenWrt BusyBox 没装 fuser/lsof，用 ps + grep + kill 兼容
	for pid in $(ps w | grep -E 'python3.*mahjong_mitm' | grep -v grep | awk '{print $1}'); do
		kill -9 "$pid" 2>/dev/null
	done
	sleep 1

	procd_open_instance
	procd_set_param command "$PYTHON" -m mahjong_mitm $args
	procd_set_param env PYTHONPATH="$PROG_DIR"
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance

	# 启动轻量 watchdog：定时 probe /healthz，连续失败则 restart 服务
	# ⚠ 必须 setsid 完全脱离父 shell，否则父 init.d 进程永不返回（阻塞 start/stop/restart）
	WATCHDOG_BIN="$PROG_DIR/watchdog.sh"
	if [ -x "$WATCHDOG_BIN" ]; then
		setsid sh "$WATCHDOG_BIN" "$tls_port" </dev/null >/dev/null 2>&1 &
		echo $! > /var/run/mahjong-mitm-watchdog.pid
	fi
}

stop_service() {
	clear_nft
	# 停止 watchdog
	if [ -f /var/run/mahjong-mitm-watchdog.pid ]; then
		kill $(cat /var/run/mahjong-mitm-watchdog.pid) 2>/dev/null || true
		rm -f /var/run/mahjong-mitm-watchdog.pid
	fi
	rm -f /var/run/mahjong-mitm-watchdog.sh
	# 兜底：把所有 python3 -m mahjong_mitm 进程清掉（procd 偶有遗漏）
	for pid in $(ps w | grep -E 'python3.*mahjong_mitm' | grep -v grep | awk '{print $1}'); do
		kill -9 "$pid" 2>/dev/null
	done
}

service_triggers() {
	procd_add_reload_trigger "mahjong-mitm"
}
